# supervisor/cloud_sql_proxy.conf
[program:SET_NAME]
command=/usr/local/bin/COMMAND
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/COMMAND.err.log
stdout_logfile=/var/log/COMMAND.out.log
user=wefitter